# infra/backup/cron/genesis-backup.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=""

17 3 * * * shellr /projects/genesis/scripts/backup.sh >> /projects/genesis/infra/backup/logs/cron.log 2>&1
